# rc.d/mosdns
#!/bin/sh

# PROVIDE: mosdns
# REQUIRE: SERVERS
# KEYWORD: shutdown

. /etc/rc.subr

name="mosdns"
rcvar="mosdns_enable"

pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"
command="/usr/sbin/daemon"
mosdns_env="IS_DAEMON=1"

daemon_args="-P ${pidfile} -r -f -o ${logfile} -t ${name}:daemon"
command_args="${daemon_args} env ${mosdns_env} /usr/local/bin/mosdns start \
    --as-service \
    -d /usr/local/etc/mosdns/"

load_rc_config $name
: ${mosdns_enable:="NO"}

start_precmd="mosdns_prestart"
mosdns_prestart()
{
    if [ ! -f /usr/local/etc/mosdns/config.yaml ]; then
        echo "Error: Config file /usr/local/etc/mosdns/config.yaml not found!"
        return 1
    fi
    [ -d /var/log ] || mkdir -p /var/log
    touch ${logfile}
    if id mosdns >/dev/null 2>&1; then
        chown mosdns:mosdns ${logfile}
    fi
    chmod 640 ${logfile}
    return 0
}

status_cmd="mosdns_status"
mosdns_status()
{
    if [ -f ${pidfile} ]; then
        pid=$(cat ${pidfile})
        if pgrep -F ${pidfile} >/dev/null 2>&1; then
            echo "mosdns is running as PID ${pid}."
            return 0
        else
            echo "mosdns is not running (stale PID file)."
            return 1
        fi
    else
        echo "mosdns is not running."
        return 1
    fi
}

stop_cmd="mosdns_stop"
mosdns_stop()
{
    if [ -f ${pidfile} ]; then
        kill $(cat ${pidfile}) && echo "mosdns stopped."
        rm -f ${pidfile}
    else
        echo "mosdns is not running."
    fi
}

run_rc_command "$1"
